# runit/gpm-runit/run
#!/bin/sh
exec 2>&1
. /etc/runit/sysconfig/mouse
exec gpm -D -m ${MDEVICE} -t ${PROTOCOL} ${GPMOPTS} 1> /dev/null 2> /dev/null
